# system/mnt-synology-backup.automount
[Unit]
Description=Automount Synology - Backup

[Automount]
Where=/mnt/synology/backup
TimeoutIdleSec=0

[Install]
WantedBy=multi-user.target
