# system/mnt-synology-music.automount
[Unit]
Description=Automount Synology - Music

[Automount]
Where=/mnt/synology/music
TimeoutIdleSec=0

[Install]
WantedBy=multi-user.target
